# files/etc/systemd/system/auto-hibernate.service
[Unit]
Description=Hibernate

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/auto-hibernate-and-resume
