# files/etc/systemd/system/auto-shutdown.service
[Unit]
Description=Shutdown and Resume

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/auto-shutdown-and-resume
