# files/etc/systemd/system/auto-resume.service
[Unit]
Description=Does nothing

[Service]
Type=oneshot
ExecStart=/bin/echo Resumed
