# files/etc/systemd/system/snapraid-maintenance.service
[Unit]
Description=SnapRAID maintenance
 
[Service]
User=root
ExecStart=/bin/bash -c "snapraid sync && snapraid scrub"
