# files/etc/systemd/system/lvm2-pvscan@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/lvm pvscan --cache --partial --activate ay %i
